# contrib/mqttlogger-restart.service
[Unit]
Description=Restart mqttlogger service (user scope)
[Service]
Type=oneshot
ExecStart=systemctl --user restart mqttlogger.service

